# deploy/agent/rtime-status-agent.timer
[Unit]
Description=Run RTime Status Board metrics agent every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
RandomizedDelaySec=10s
AccuracySec=5s
Unit=rtime-status-agent.service

[Install]
WantedBy=timers.target
